# luci-app-lucky/root/etc/uci-defaults/99_x_luckyarch.sh
#!/bin/sh

chmod +x /usr/bin/luckyarch
exit 0